# scripts/systemd/stellwerk-preprod-autodeploy.service
[Unit]
Description=Stellwerk preprod auto-deploy (pull + docker compose)
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot

# Configure these via /etc/stellwerk-preprod-autodeploy.env (recommended)
Environment=REPO_DIR=/opt/goals
Environment=DEPLOY_FLAGS=
EnvironmentFile=-/etc/stellwerk-preprod-autodeploy.env

ExecStart=/usr/bin/env bash -lc 'cd "${REPO_DIR}" && ./scripts/preprod-autodeploy.sh'
